# sort-downloads.service
[Unit]
Description=Sort Downloads

[Service]
Type=oneshot
ExecStart=/home/santeri/git/scripts/sort-downloads.py
