# sort-downloads.timer
[Unit]
Description=Sort Downloads monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=default.target
